# build/crontab
0 0-23 * * * php protected/t4.php notify/guests
